# dot_config/autostart/xbindkeys.desktop
[Desktop Entry]
Name=xbindkeys
Comment=Start xbindkeys
Exec=xbindkeys
Type=Application
